# backend/.ops/sit/supervisor.d/horizon.conf
[program:horizon]
command=php artisan horizon
directory=/var/www/html
user=nginx
autostart=true
autorestart=true
priority=6
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/horizon.log
