# .ops/etc/supervisor/conf.d/scheduler.conf
[program:scheduler]
command=./schedule.sh
directory=/srv/www
user=nginx
environment=HOME="/srv/www",USER="nginx"
